# forecasting-classification/crontab
@reboot cd ~/../app && /usr/local/bin/python forecasting.py >> /var/log/cron.log 2>&1
0 */8 * * * cd ~/../app && /usr/local/bin/python forecasting.py >> /var/log/cron.log 2>&1
